# .local/bin/setup/screenlock/screenlock@.service
[Unit]
Description=Starts i3lock for user %i at suspend time
Before=sleep.target

[Service]
User=%i
Type=forking
Environment=DISPLAY=:0
ExecStartPre=
ExecStart=/home/%i/.local/bin/lock

[Install]
WantedBy=sleep.target
